# deploy/chron0-tui.service
[Unit]
Description=chron0-tui — Jon Marien's TUI portfolio (SSH)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tui
Group=tui
WorkingDirectory=/opt/chron0-tui
ExecStart=/home/tui/.bun/bin/bun src/index.ts
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=2222
Environment=HOST_KEY_PATH=/home/tui/.ssh/host_rsa
Environment=MAX_CONNECTIONS=50
Environment=MAX_PER_IP=5
Restart=on-failure
RestartSec=3s
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/tui/.ssh
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
